# runpod/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:comfyui]
command=/app/start_comfy.sh
directory=/app/ComfyUI
environment=PYTHONPATH="/app/ComfyUI",PYTHONUNBUFFERED="1",PYTHONIOENCODING="utf-8"
autostart=true
autorestart=true
startsecs=15
stdout_logfile=/var/log/comfyui.log
stderr_logfile=/var/log/comfyui_err.log
priority=20

[program:backend]
command=python3 -u /app/backend/server.py
directory=/app/backend
environment=PYTHONPATH="/app/backend:/app/ComfyUI",PYTHONUNBUFFERED="1",CORS_ORIGINS="*",COMFY_VIEW_BASE="/comfy",COMFY_URL="http://127.0.0.1:8199"
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/backend.log
stderr_logfile=/var/log/backend_err.log
priority=30

[program:jupyter]
command=python3 -m jupyterlab --ServerApp.ip=0.0.0.0 --ServerApp.port=8888 --ServerApp.open_browser=False --ServerApp.allow_root=True --ServerApp.token='' --ServerApp.password='' --ServerApp.root_dir=/workspace --ServerApp.preferred_dir=/workspace --ServerApp.allow_origin=* --ServerApp.trust_xheaders=True --ServerApp.disable_check_xsrf=True --ServerApp.terminals_enabled=True
directory=/workspace
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/jupyter.log
stderr_logfile=/var/log/jupyter_err.log
priority=35

[program:nginx]
command=nginx
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx_access.log
stderr_logfile=/var/log/nginx_error.log
priority=40

